# scripts/git-indicator.cron
#!/bin/sh
# To be added in cron file

pid=`git-indicator pid`

[ "$pid" = "" ] && exit 0

rss=`awk -F " " '/VmRSS/{ printf $2 }' /proc/$pid/status`

[ $rss -ge 60000 ] && git-indicator restart
