# scripts/git-indicator.cron
#!/bin/sh

while true
do
    pid=`git-indicator pid`

    [ -f "/proc/$pid/status" ] && exit 0

    rss=`awk -F " " '/VmRSS/{ printf $2 }' /proc/$pid/status`

    [ $rss -ge 60000 ] && git-indicator restart

    sleep 5
done
